# cloudflared/cloudflared-api.service.example
[Unit]
Description=Cloudflare Tunnel for api.clipscut.pro
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
ExecStart=/usr/bin/cloudflared tunnel --config /etc/cloudflared/api-clipscut.yml --no-autoupdate run
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
